# imageroot/systemd/user/freshclam.timer
[Unit]
Description=Timer of freshclam service

[Timer]
OnUnitInactiveSec=60 minutes

[Install]
WantedBy=timers.target
