# imageroot/systemd/user/spam-expunge.timer
[Unit]
Description=Timer of spam-expunge procedure

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
